# docker/s6-rc.d/integrasjonspunkt/run
#!/command/with-contenv bash

cd ${APP_DIR}
# Korrigerer integrasjonspunktets port dersom porten endres i Azure Web App
if [ -n "${PORT+set}" ]; then
    export SERVER_PORT=$PORT
fi
exec s6-setuidgid java java ${JAVA_OPTS} -jar app.jar
