# docker/s6-rc.d/integrasjonspunkt/run
#!/command/with-contenv bash

cd ${APP_DIR}
exec s6-setuidgid java java ${JAVA_OPTS} -jar app.jar
