# docker/s6-rc.d/integrasjonspunkt/run
#!/command/with-contenv bash

cd ${APP_DIR}
# Korrigerer integrasjonspunktets port dersom porten endres i Azure Web App-konfig
if [ -n "${PORT+set}" ]; then
    export SERVER_PORT=$PORT
fi
unset JAVA_TOOL_OPTIONS
runAsUser=www-data
if [ "$SERVER_PORT" -lte "1024" ]; then
    # Må kjøre som root med SERVER_PORT < 1024
    runAsUser=root
fi
exec s6-setuidgid $runAsUser java ${JAVA_OPTS} -jar app.jar
