# docker/services.d/integrasjonspunkt/run
#!/usr/bin/with-contenv bash

s6-setuidgid java ${JAVA_HOME}/bin/java ${JAVA_OPTS} -Dspring.profiles.active=${APP_ENV} -jar ${APP_DIR}/app.jar
