# docker/services.d/integrasjonspunkt/run
#!/command/execlineb -P
with-contenv
s6-setuidgid java 
java ${JAVA_OPTS} -Dspring.profiles.active=${APP_ENV} -jar ${APP_DIR}/app.jar
